# base-files/files/etc/crontabs/root

* * * * * /bin/sh /etc/config/adg/check_agh.sh
