# base-files/files/etc/rc.d/S99ledoff
#!/bin/sh /etc/rc.common
START=99

start() {
    sleep 2
    for led in /sys/class/leds/*; do
        echo 0 > "$led/brightness"
    done
}

